# backend/services/agricontrol-filter.service
[Unit]
Description=AgriControl Kalman Filter Engine
After=network.target mosquitto.service influxdb.service agricontrol-mqtt.service
Wants=mosquitto.service influxdb.service

[Service]
Type=simple
User=sguser
WorkingDirectory=/home/sguser/iot-backend
ExecStart=/home/sguser/iot-backend/venv/bin/python3 -u filter_engine.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/sguser/iot-backend/filter.log
StandardError=append:/home/sguser/iot-backend/filter.log

[Install]
WantedBy=multi-user.target
